# systemd/all-youtube-summary.service
[Unit]
Description=All YouTube summary/publish worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/all-youtube
EnvironmentFile=/home/ubuntu/all-youtube/.env
ExecStart=/home/ubuntu/all-youtube/.venv/bin/python -m app.worker_summary
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target